# scripts/voice-hermes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voice-hermes.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>-m</string>
        <string>voice_hermes</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/riandy/development/hermes/voice-hermes</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/voice-hermes.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/voice-hermes.err</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
